# terraform/files/gnbsim/rc.local
#!/bin/bash

sudo ip link add ran link enp6s0 type macvlan mode bridge
sudo ip link set dev ran up
